# debian/mcsauna.upstart
description "daemon to report back hot memcached keys in a graphite-friendly format"

start on networking
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid Debian-mcsauna

exec /usr/bin/mcsauna -c /etc/mcsauna/mcsauna.conf
